# host/backup.service
[Unit]
Description=Remote backup to encrypted device

[Service]
Type=simple
User=root
ExecStart=/usr/bin/backup.sh
